# ansible/roles/lolcat/files/catinder-api-update.path
[Path]
PathChanged=/home/lolcat/catinder-api.jar

[Install]
WantedBy=multi.user.target
