# ansible/roles/lolcat/files/catinder-api-update.service
[Unit]
Description=catinder-api restart on update

[Service]
ExecStart=/bin/systemctl restart catinder-api
