# ansible/roles/lolcat/files/proximity-aggregator-update.path
[Path]
PathChanged=/home/lolcat/proximity-aggregator.jar

[Install]
WantedBy=multi.user.target
