# ansible/roles/lolcat/files/proximity-aggregator-update.service
[Unit]
Description=proximity-aggregator restart on update

[Service]
ExecStart=/bin/systemctl restart proximity-aggregator
