# ansible/roles/lolcat/files/simple-api-update.path
[Path]
PathChanged=/home/lolcat/simple-api.jar

[Install]
WantedBy=multi.user.target
