# ansible/roles/lolcat/files/stream-update.service
[Unit]
Description=simple-api restart on update

[Service]
ExecStart=/bin/systemctl restart simple-api
